# systemd/powerloss-shutdown.service
[Unit]
Description = [pipower] Power off system

[Service]
Type = simple
ExecStart = /bin/systemctl poweroff
