# systemd/powerloss-shutdown.timer
[Unit]
Description = [pipower] Power loss timer

[Timer]
OnActiveSec = 30

[Install]
WantedBy = multi-user.target
